# systemd/audioscribe-api.service
[Unit]
Description=Audio Scribe API
After=network.target postgresql.service redis-server.service audioscribe-minio.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=talinoserver4
Group=talinoserver4
WorkingDirectory=/home/talinoserver4/Documents/audio-scribe/backend
EnvironmentFile=/home/talinoserver4/Documents/audio-scribe/backend/.env
ExecStart=/home/talinoserver4/Documents/audio-scribe/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
